# deploy/mc-backup.service
[Unit]
Description=Mission Control nightly SQLite backup

[Service]
Type=oneshot
User=mc
Group=mc
EnvironmentFile=/srv/mission-control/.env
ExecStart=/srv/mission-control/deploy/mc-backup.sh
